# supervisor/dockerd.conf
[program:dockerd]
environment=BUILDKIT_STEP_LOG_MAX_SIZE=95971520,BUILDKIT_STEP_LOG_MAX_SPEED=7048576
command=/usr/bin/dockerd
autostart=true
autorestart=true
stderr_logfile=/var/log/dockerd.err.log
stdout_logfile=/var/log/dockerd.out.log
